# sno-install-container/installer/iptables.service.tmpl
[Unit]
Description=Setup openshift iptables
[Service]
ExecStart=/usr/local/bin/ocp-iptables.sh
Type=oneshot
RemainAfterExit=true
[Install]
WantedBy=multi-user.target

